# nextcloud/nextcloud-sync.timer
[Unit]
Description=Run Nextcloud Sync Periodically

[Timer]
OnBootSec=5m
OnUnitActiveSec=30m
RandomizedDelaySec=2m
Unit=nextcloud-sync.service

[Install]
WantedBy=timers.target
